# scripts/crystal-shell.service
[Unit]
Description=Crystal Shell UI (bar, dock, overlays)
Documentation=https://github.com/fluid-crystal/crystal-shell
PartOf=graphical-session.target
After=graphical-session.target
# Crash-loop guard: if the shell dies more than 5 times in 30s (e.g. a broken
# build or a config that crashes on startup), stop retrying instead of respawning
# forever in a flickering loop — drop to a TTY and fix it. Re-arm with
# `systemctl --user reset-failed crystal-shell.service`.
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/crystal-shell-ui
# Respawn on crash (segfault / non-zero exit / signal). A clean exit is NOT
# restarted, so quitting on purpose stays quit.
Restart=on-failure
RestartSec=1
# Keep restarts snappy: SIGKILL the cgroup quickly if it doesn't exit on SIGTERM.
TimeoutStopSec=3
Slice=session-graphical.slice

# NOTE: intentionally NO [Install]/WantedBy=graphical-session.target. The Crystal
# session and a plain uwsm-managed Hyprland session are indistinguishable at the
# systemd-target level (both are compositor "Hyprland" reaching the same generic
# graphical-session.target), so binding here would start the Crystal UI in ANY
# Hyprland session — leaking the bar/dock into sessions that aren't ours. Instead
# this service is started explicitly from the Crystal Hyprland config
# (config/hypr/hyprland.lua → `systemctl --user start crystal-shell.service`),
# which only loads in the Crystal session. PartOf= still stops it cleanly on logout.
